# mongo.conf
[program:redis]
command=mongo-entrypoint.sh mongod
