# mongo.conf
[program:mongodb]
command=mongo-entrypoint.sh mongod --storageEngine=mmapv1
